# recipes-core/systemd/rtc-sanity-check/rtc-sanity-check.service
[Unit]
Description=rtc sanity check

[Service]
ExecStart=/usr/sbin/rtc-sanity-check.sh
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=basic.target
